# .docker/crontab
0,15,30,45 9-20 * * * /opt/console -sport=football > /proc/1/fd/1 2> /proc/1/fd/2
