# .docker/crontab
14,29,44,59 9-20 * * * /opt/console -sport=football > /proc/1/fd/1 2> /proc/1/fd/2
